# my-packages/sync-zeronet/sync-zeronet@.timer
[Unit]
Description=Sync zeronet users.json Minutely

[Timer]
OnCalendar=minutely
Persistent=true

[Install]
WantedBy=multi-user.target
