# my-packages/my-notmuch-utils/delete-mails.service
[Unit]
Description=Delete notmuch deleted and spam mails Service
AssertFileIsExecutable=/usr/local/bin/delete-mails.sh

[Service]
Type=oneshot
ExecStart=/usr/local/bin/delete-mails.sh
IOSchedulingClass=idle
IOSchedulingPriority=7
