# my-packages/my-notmuch-utils/sync-mails.timer
[Unit]
Description=Synchronize mails Timer. 5m

[Timer]
OnStartupSec=10m
OnUnitActiveSec=5m

[Install]
WantedBy=timers.target
